# scripts/influxdb_uploader/mosogep.service
[Unit]
Description=Mosogep Collector
After=network.target

[Service]
ExecStart=/usr/bin/python3 /opt/save.py
#KillMode=process
Restart=on-failure
